# stats/deploy/magies-shell-release-sync.service
[Unit]
Description=Sync MagiesTerminal GitHub releases to shell.magies.top mirror
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Environment=RELEASE_DEST=/var/www/shell.magies.top/releases
# Optional: set GITHUB_TOKEN in this file for higher API rate limits
# EnvironmentFile=-/etc/magies-shell/release-sync.env
ExecStart=/opt/magies-shell-site/scripts/sync-github-releases.sh
Nice=10
